# A2CS/code/autostart/iot.desktop
[Desktop Entry]
Type=Application
Name= IOT
Terminal=true
Exec=lxterminal -e /home/pi/Documents/./IOTscript.sh